# config/systemd/sdr-pi-op25.service
[Unit]
Description=sdr-pi OP25 P25 scanner
After=network-online.target sdr-pi-performance.service
Wants=network-online.target
Requires=sdr-pi-performance.service

[Service]
Type=simple
EnvironmentFile=/etc/sdr-pi/sdr-pi.conf
ExecStart=/usr/local/bin/sdr-pi-op25-wrapper
Restart=on-failure
RestartSec=5
User=sdr
Group=plugdev

# SDR performance tuning — pin to dedicated core, elevate priority,
# lock sample buffers in RAM. Core 0 is reserved for USB IRQs.
Nice=-10
CPUAffinity=3
LimitRTPRIO=50
LimitMEMLOCK=infinity
LimitNOFILE=8192

[Install]
WantedBy=multi-user.target
