# config/systemd/sdr-pi-dump1090.service
[Unit]
Description=sdr-pi dump1090 ADS-B receiver
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/sdr-pi/sdr-pi.conf
ExecStart=/usr/local/bin/sdr-pi-dump1090-wrapper
Restart=on-failure
RestartSec=5
User=sdr
Group=plugdev

[Install]
WantedBy=multi-user.target
